# linux/x/sysdunit/bar.service
[Service]
User=bar
ExecStart=/usr/local/bin/bar.py

[Install]
WantedBy=multi-user.target
